# deploy/botmarket-api.service
# systemd service for BotMarketplace API (Fastify + Bot Worker)
# Install: cp deploy/botmarket-api.service /etc/systemd/system/
# Enable:  systemctl daemon-reload && systemctl enable botmarket-api
# Start:   systemctl start botmarket-api
# Logs:    journalctl -u botmarket-api -f

[Unit]
Description=BotMarketplace API
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=botmarket
WorkingDirectory=/opt/-botmarketplace-site/apps/api
EnvironmentFile=/opt/-botmarketplace-site/.env

# Safety net: kill any orphan process holding port 4000 before starting.
# Leading "-" means failure (e.g. no process on the port) is ignored.
# Addresses recurring orphan-process issue where previous node/next workers
# survive `systemctl restart` with PPID=1 and block port binding.
ExecStartPre=-/usr/bin/fuser -k 4000/tcp

ExecStart=/usr/bin/node dist/server.js

# Graceful shutdown: SIGTERM to main process, SIGKILL to any lingering
# children after TimeoutStopSec. server.ts installs SIGTERM handler that
# closes Fastify + Prisma within 30s. SendSIGKILL guarantees cleanup.
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=35s
SendSIGKILL=yes

Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=botmarket-api

# Ensure node_modules/.bin is reachable (for prisma etc.)
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin:/usr/local/share/pnpm

# Hardening (docs/30 Phase 9). MemoryDenyWriteExecute is intentionally
# omitted — V8 JIT requires writable+executable pages.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ReadWritePaths=/opt/-botmarketplace-site /var/log

[Install]
WantedBy=multi-user.target
